# _pi_remote/etc_systemd_system_camstream.service
[Unit]
Description=Camera stream (GPIO17 GND → app по домашней Wi-Fi)
After=multi-user.target network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=pavel
Group=pavel
WorkingDirectory=/home/pavel/projects/Mobile_Raspberry_5-
Environment=CAMSTREAM_BOOT_GPIO_BCM=17
SupplementaryGroups=gpio
Environment=BOOT_GPIO_STABLE_MS=400
Environment=BOOT_GPIO_WAIT_GROUND_SEC=60
Environment=CAMSTREAM_HOME_WIFI_CONN=RearWindow
ExecStart=/bin/sh -lc 'exec /home/pavel/projects/Mobile_Raspberry_5-/scripts/camstream-gpio-wifi-fallback.sh send --video-mode udp_h264 --host 192.168.1.16 --port 5000 --stream-preset realtime'

Restart=no

[Install]
WantedBy=multi-user.target
